# deploy/systemd/quant-kernel.service
[Unit]
Description=Reco Trading QuantKernel
After=network-online.target redis.service postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=trading
Group=trading
WorkingDirectory=/opt/reco-trading
EnvironmentFile=-/opt/reco-trading/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/env bash -lc 'cd /opt/reco-trading && python main.py'
Restart=always
RestartSec=5
TimeoutStopSec=30
LimitNOFILE=262144

[Install]
WantedBy=multi-user.target
