# reco_trading/system/reco_trading.service
[Unit]
Description=Reco-Trading Bot Daemon
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
WorkingDirectory=%h/Downloads/trading/reco-trading
ExecStart=/usr/bin/python -m reco_trading trade --config config.json
Restart=on-failure
RestartSec=10
User=%u
Environment=PYTHONUNBUFFERED=1

StandardOutput=journal
StandardError=journal
SyslogIdentifier=reco_trading

[Install]
WantedBy=multi-user.target
